# ops/systemd/vidgen-control-center-ingest.service
[Unit]
Description=Vidgen Control Center ingest snapshot job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=clawuser
Group=clawuser
WorkingDirectory=/docker/openclaw-jnqf/data/repos/vidgen/apps/control-center
EnvironmentFile=-/etc/default/vidgen-control-center
ExecStart=/usr/bin/env npm run ingest:snapshot
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/docker/openclaw-jnqf/data/repos/vidgen/apps/control-center/data
